# roles/common/templates/etc/profile.d/ros-activate.sh.j2
source /opt/ros/{{ ros_release }}/setup.bash
