# roles/controls-accelerator/files/etc/profile.d/controls-accelerator.sh
alias help="/bin/cat /opt/controls-accelerator/HELP"
